# ratatunel5000.conf
[program:ratatunel5000]
command=ssh -i /home/pi/.ssh/id_rsa -N  -o "ServerAliveInterval 20" -o ExitOnForwardFailure=yes -R 5000:localhost:5000    root@207.246.118.54
directory=/home/pi
autostart=true
user=pi
stdout_logfile=/home/pi/commandcenter/data/logs/ratatunel5000.log
redirect_stderr=true
numprocs=1
